# ops/templates/systemd/prophet_labs_scheduler.service
[Unit]
Description=Prophet Labs scheduler service
After=network.target

[Service]
EnvironmentFile=/path/to/prophet-labs/.env
WorkingDirectory=/path/to/prophet-labs
ExecStart=/path/to/prophet-labs/.venv/bin/python -m prophet_labs.jobs.scheduler
Restart=on-failure
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
